# rc.local
#!/bin/bash

/home/bvk/config/99-macbook-fixes.sh bootup
exit 0
